# service-monitor.service
[Unit]
Description=Headless service monitor (live status page for Claude services + routines)
After=default.target

[Service]
Type=simple
# Bind address. 127.0.0.1 = localhost-only (view via SSH tunnel). Set to a
# specific interface IP (ideally a VPN/mesh address) to expose it there. Avoid
# 0.0.0.0 unless a firewall scopes the port. IP_FREEBIND handles interfaces that
# appear after boot.
Environment=MONITOR_BIND=127.0.0.1
Environment=MONITOR_PORT=8787
ExecStart=/usr/bin/python3 %h/routines/monitor/monitor.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
